# src/assets/systemd.service
[Unit]
Description=Ourcraft
After=network.target

[Service]
Type=forking
User=${USER}
Group=${GROUP}
Restart=on-failure
RestartSec=3
ExecStart=/usr/bin/screen -L -dmS ourcraft  ${ROOT}/ourcraft start

[Install]
WantedBy=multi-user.target
